# cfs-ui/config/autostart/autotightvnc.desktop
[Desktop Entry]
Type=Application
Name=TightVNC
Exec=vncserver :0
StartupNotify=false
