# cfs-ui/config/autostart/tightvnc.desktop
[Desktop Entry]
Type=Application
Name=TightVNC
Exec=lxterminal --command "vncserver :1"
StartupNotify=false
Categories=Utility;
Icon=preferences-desktop-remote-desktop
